# bootstrap/rc.local
#!/bin/sh

# Set serial number (MAC address) as default hostname
echo "Set hostname ..."
if [ ! -f /etc/hostname ]; then
        echo "iSmartAlarm" > /etc/hostname
fi
/bin/hostname -F /etc/hostname

ifconfig lo up

echo 45 > /proc/sys/net/ipv4/tcp_keepalive_time
echo 3 > /proc/sys/net/ipv4/tcp_keepalive_probes
echo 5 > /proc/sys/net/ipv4/tcp_keepalive_intvl

if [ -x "/etc/fang_hacks.sh" ]; then
    /etc/fang_hacks.sh
fi
